# config/com.calos.cal-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service Label (unique identifier) -->
    <key>Label</key>
    <string>com.calos.cal-daemon</string>

    <!-- Program to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/matthewmauer/Desktop/CALOS_ROOT/agent-router/bin/cal-daemon</string>
    </array>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/matthewmauer/Desktop/CALOS_ROOT/agent-router</string>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>CAL_INTERVAL</key>
        <string>1h</string>
        <key>CAL_AUTO_PUBLISH</key>
        <string>false</string>
        <key>CAL_AUTO_SIGN</key>
        <string>true</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- Auto-start on boot -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Keep alive (restart if crashes) -->
    <key>KeepAlive</key>
    <true/>

    <!-- Standard output log -->
    <key>StandardOutPath</key>
    <string>/Users/matthewmauer/Desktop/CALOS_ROOT/agent-router/logs/cal-daemon-out.log</string>

    <!-- Standard error log -->
    <key>StandardErrorPath</key>
    <string>/Users/matthewmauer/Desktop/CALOS_ROOT/agent-router/logs/cal-daemon-error.log</string>

    <!-- Restart on failure -->
    <key>ThrottleInterval</key>
    <integer>60</integer>

    <!-- Resource limits -->
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>1024</integer>
    </dict>

    <key>HardResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>2048</integer>
    </dict>
</dict>
</plist>
